# systemd/superhot-tty@.service
[Unit]
Description=SuperHot TTY window manager v0.5 on %I (privilege separation)
Conflicts=rescue.service
After=systemd-logind.service systemd-vconsole-setup.service
Before=getty.target
Wants=systemd-logind.service
After=systemd-user-sessions.service plymouth-quit.service

# Disable default agetty — superhot-tty имеет собственный login screen.
Conflicts=getty@%i.service

# Crash loop protection: если сервис упал 5 раз за 60 секунд,
# systemd перестанет перезапускать. Дополнительно WM сам откатывается
# на getty при обнаружении crash loop (см. check_crash_loop в main.rs).
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=idle
ExecStart=/usr/local/bin/superhot-tty
Restart=always
RestartSec=2

StandardInput=tty
StandardOutput=journal
StandardError=journal
TTYPath=/dev/%I
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

# Запуск от root — нужен только на старте:
#   1. Открыть DRM master fd
#   2. Открыть /dev/input/event*
#   3. fork() — ребёнок drop к "superhot-tty" (login screen)
#      родитель остаётся root для PAM auth, затем drop к вошедшему пользователю
# После успешного входа процесс работает от имени пользователя.
#
# ВАЖНО: НЕ используем PAMName=login — systemd пытается создать login session
# через pam_systemd (Varlink), что ломает запуск (InvalidParameter error).
# WM сам делает PAM auth через FFI к libpam в privilege-separated child.
User=root
WorkingDirectory=/

# Hardening (совместимо с DRM master requirements).
ProtectSystem=full
ProtectHome=false
PrivateTmp=true
PrivateDevices=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

Environment=DISPLAY=:1
Environment=XDG_SESSION_TYPE=desktop
Environment=XDG_CURRENT_DESKTOP=superhot
Environment=LANG=en_US.UTF-8
Environment=RUST_LOG=info
Environment=SUPERHOT_TTY=1

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
DefaultInstance=tty1
