# systemd/superhot-tty@.service
[Unit]
Description=SuperHot TTY window manager v0.2 on %I
Conflicts=rescue.service
After=systemd-logind.service systemd-vconsole-setup.service
Before=getty.target
Wants=systemd-logind.service pipewire.service pipewire-pulse.service wireplumber.service
After=systemd-user-sessions.service plymouth-quit.service

# Disable default agetty.
Conflicts=getty@%i.service

[Service]
Type=idle
ExecStart=/usr/local/bin/superhot-tty
Restart=always
RestartSec=2

StandardInput=tty
TTYPath=/dev/%I
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

User=root
WorkingDirectory=/root
PAMName=login

ProtectSystem=full
ProtectHome=false
Environment=DISPLAY=:1
Environment=XDG_RUNTIME_DIR=/run
Environment=XDG_SESSION_TYPE=desktop
Environment=XDG_CURRENT_DESKTOP=superhot
Environment=LANG=en_US.UTF-8
Environment=RUST_LOG=info
Environment=SUPERHOT_TTY=1

# PipeWire socket dir (если запускаем pipewire от root).
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
DefaultInstance=tty1
